# images/cryptpad/etc/s6-overlay/s6-rc.d/nginx/run
#!/bin/sh

exec 2>&1 # redirects stderr to stdout

exec nginx -g "daemon off;" -c /etc/nginx/nginx-final.conf
